# server/netbox.service
[Unit]
Description='Netbox Service'
Requires=docker.service
After=docker.service

[Service]
WorkingDirectory=/home/lab/dev/st2110-automation-bench/server/netbox/
ExecStart=/usr/bin/env /usr/local/bin/docker-compose up
ExecStop=/usr/bin/env /usr/local/bin/docker-compose  stop

[Install]
WantedBy=multi-user.target
